# roles/prometheus/templates/prometheus.service.j2
[Unit]
Description=Prometheus

[Service]
ExecStart=/usr/bin/prometheus --config.file=/etc/prometheus.yml
Restart=Always
