# roles/filebrowser/templates/filebrowser.service.j2
[Unit]
Description=FileBrowser

[Service]
WorkingDirectory=/media/data/archive

ExecStart=/usr/bin/filebrowser -p 40197
Restart=Always

[Install]
Wants=network.target
